# comet@.service
[Unit]
Description=RT2000BT BLE->MQTT bridge (%i)
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/rt2000bt
EnvironmentFile=/opt/rt2000bt/.env
ExecStart=/opt/rt2000bt/.venv/bin/python /opt/rt2000bt/comet.py
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
